# filechooser-portal/io.elementary.files4.xdg-desktop-portal.service.in
[Unit]
Description=Portal service (Pantheon Files)

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.elementary.files4
ExecStart=@libexecdir@/io.elementary.files4.xdg-desktop-portal
